# technics/K8062dllSourceCode/K8062ForLinux/DmxControl_V2/Dmxd/10-k8062.rules
KERNEL=="1-3.2", SUBSYSTEMS=="usb", ATTRS{product}=="USB K8062", RUN+="/root/dmxd"